# debian/dwarf.upstart
description "Dwarf server"
author "Juerg Haefliger <juerg.haefliger@hp.com>"

start on (local-filesystem and net-device-up IFACE=lo) or (started libvirt-bin)

respawn

exec su -c "dwarf" dwarf
